# templates/multi-node-system/nodes/pi/systemd/internal-relay.service.template
[Unit]
Description=Example internal relay service
After=network-online.target

[Service]
Type=simple
User=<service-user>
WorkingDirectory=<path-to-working-copy>/nodes/pi/<service-name>
ExecStart=<path-to-working-copy>/nodes/pi/<service-name>/run.sh
Restart=on-failure

[Install]
WantedBy=multi-user.target
